# private_dot_config/exact_systemd/exact_user/exact_org.freedesktop.IBus.session.GNOME.service.d/override.conf
[Unit]
Description=
Description=IBus Daemon for GNOME (overrided by Fcitx5)

[Service]
Environment="LANG=zh_CN.UTF-8"
ExecStart=
ExecStart=/usr/bin/fcitx5 --keep -r
